# crontab
*	*	*	*	*	cd /app && python manage.py rancher_sync
