# crontab
*	*	*	*	*	cd /app && python manage.py rancher_sync
*	*	*	*	*	cd /app && python manage.py meta > /tmp/data.json
