# crontab
*	*	*	*	*	cd /app && python manage.py rancher_sync
*	*	*	*	*	cd /app && python manage.py admin_pw > /tmp/data.json
